# hp/honeytrap/dist/supervisord.conf
[supervisord]
nodaemon=true

[program:honeytrap]
command=/opt/honeytrap/sbin/honeytrap -D -C /opt/honeytrap/etc/honeytrap/honeytrap.conf -t 5 -u tpot -g tpot
autorestart=true

[program:ewsposter]
command=/usr/bin/python /opt/ewsposter/ews.py -c /data/ews/conf/ -m honeytrap -l 60
autorestart=true
