# hp/dionaea/dist/supervisord.conf
[supervisord]
nodaemon=true

[program:dionaea]
command=/opt/dionaea/bin/dionaea -u tpot -g tpot -c /opt/dionaea/etc/dionaea/dionaea.cfg -p /data/dionaea/dionaea.pid
autorestart=true

[program:ewsposter]
command=/usr/bin/python /opt/ewsposter/ews.py -c /data/ews/conf/ -m dionaea -l 60
autorestart=true
