# hp/elasticpot/dist/supervisord.conf
[supervisord]
nodaemon=true

[program:elasticpotpy]
directory=/opt/ElasticpotPY/
command=/usr/bin/python main.py
autorestart=true
