# hp/conpot/dist/supervisord.conf
[supervisord]
nodaemon=true

[program:conpot]
command=/usr/local/bin/conpot --template kamstrup_382 --logfile /data/conpot/log/conpot.log --config /etc/conpot/conpot.cfg
autorestart=true
