# hp/glastopf/dist/supervisord.conf
[supervisord]
nodaemon=true

[program:glastopf]
directory=/data/glastopf/
command=bash -c 'cp /opt/glastopf/glastopf.cfg /data/glastopf/ && exec glastopf-runner'
autorestart=true

[program:ewsposter]
command=/usr/bin/python /opt/ewsposter/ews.py -c /data/ews/conf/ -m glastopfv3 -l 60
autorestart=true
